# com.tradingagents.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tradingagents.web</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/tmp/TradingAgents/web/app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/tmp/TradingAgents</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/TradingAgents/output/web.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/TradingAgents/output/web.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/tmp/TradingAgents</string>
        <key>PORT</key>
        <string>8080</string>
    </dict>
</dict>
</plist>
